# templates/systemd/agent-self-health-full.timer
[Unit]
Description=Run OpenClaw Agent Self-Health Full Mode Hourly

[Timer]
OnBootSec=5m
OnUnitActiveSec=1h
Persistent=true
Unit=agent-self-health-full.service

[Install]
WantedBy=timers.target
